# linux/2024thinkpadx1gen8/.config/systemd/user/gnome-keyring.service
[Unit]
Description=Start gnome-keyring-daemon

[Service]
Type=simple
StandardOutput=journal
ExecStart=/home/ohthehugemanatee/bin/start-gnome-keyring.sh

[Install]
WantedBy=default.target
